# tools/lobot-collector.service
# Install:
#   sudo cp /opt/Lobot/tools/lobot-collector.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now lobot-collector
#
# Venv setup (one-time):
#   python3 -m venv /opt/Lobot/tools/lobot_collector/.venv
#   /opt/Lobot/tools/lobot_collector/.venv/bin/pip install -r /opt/Lobot/tools/lobot_collector/requirements-collector.txt
#
# Useful commands:
#   sudo systemctl status lobot-collector
#   sudo journalctl -u lobot-collector -f
#   curl -N http://localhost:9095/api/events          # verify SSE stream
#   curl -s http://localhost:9095/api/state | python3 -m json.tool

[Unit]
Description=Lobot Cluster Data Collector
After=network.target

[Service]
Type=simple
Environment=LOBOT_CLUSTER_DIR=/opt/Lobot
Environment=LOBOT_TOOLS_DIR=/opt/Lobot/tools
WorkingDirectory=/opt/Lobot/tools
ExecStart=/opt/Lobot/tools/lobot-collector.sh
User=croot
Restart=on-failure
RestartSec=15s

# Give the process time to send the shutdown email on SIGTERM
TimeoutStopSec=30s

# stdout/stderr go to journald (view with: journalctl -u lobot-collector -f)
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
