# tools/resource-collector.service
# Install:
#   sudo cp /opt/Lobot/tools/resource-collector.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now resource-collector
#
# Useful commands:
#   sudo systemctl status resource-collector
#   sudo journalctl -u resource-collector -f

[Unit]
Description=Lobot Resource Collector
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/Lobot
ExecStart=/usr/bin/python3 /opt/Lobot/resource_collector-withmail.py
User=croot
Restart=on-failure
RestartSec=15s

# Give the script time to send the shutdown email on SIGTERM
TimeoutStopSec=30s

# stdout/stderr go to journald (view with: journalctl -u resource-collector -f)
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
